# units/override/application.service
[Unit]
Description=Test application

[Service]
ExecStart=/bin/sleep infinity