# units/override/override-application.service
[Service]
ExecStart=/srv/demo/echo_script.sh