# units/service/application.service
[Unit]
Description=Test application

[Service]
ExecStart=/srv/demo/echo_script.sh